# qwen-agent.service
[Unit]
Description=Qwen Agent Container
After=network.target
Requires=qwen-network.service

[Service]
Type=simple
ExecStart=/usr/bin/podman start qwen-agent
ExecStop=/usr/bin/podman stop qwen-agent
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target